# slimlock.service
[Unit]
Description=Lock X session using slimlock

[Service]
User=madhatter
ExecStart=/usr/bin/slimlock

[Install]
WantedBy=sleep.target
